# crontab/per_minute
*/1 *    * * *   root    ps -ef | grep main.py | grep -v grep || tmux send -t spider.0 "python3 main.py" ENTER